# com.georgeolaru.limpet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!--
    LaunchAgent for limpet (per-user, runs at login).
    NOTE: install.sh automatically generates a copy with the correct paths in
    ~/Library/LaunchAgents/. This file is a reference / manual fallback.
    If you use it manually, replace /Users/georgeolaru with your home.
  -->
  <key>Label</key>
  <string>com.georgeolaru.limpet</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/georgeolaru/.local/bin/limpet.sh</string>
  </array>

  <!-- Starts at login and relaunches if the process dies. -->
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <!-- If it were to relaunch rapidly in a loop, wait at least 10s between starts. -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <!-- Low priority (background) so it doesn't eat CPU/battery. -->
  <key>ProcessType</key>
  <string>Background</string>
  <key>LowPriorityIO</key>
  <true/>

  <!-- The daemon's stdout/stderr (separate from the application log). -->
  <key>StandardOutPath</key>
  <string>/Users/georgeolaru/Library/Logs/limpet.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/georgeolaru/Library/Logs/limpet.err.log</string>
</dict>
</plist>
